# deploy/poller.cron
# Every 15 minutes
*/15 * * * * /usr/local/bin/python -m bin.poll_once
